# systemd/whisp.path
[Unit]
Description=Watch for Whisp file changes

[Path]
PathModified=/home/elec/whisp/whisp.fyi.php
PathModified=/home/elec/whisp/src/Server.php
PathModified=/home/elec/whisp/src/Channel.php
Unit=whisp.service

[Install]
WantedBy=multi-user.target
